# cloud-init-local.service.d/10-wait-for-net-device.conf
[Unit]
Requires=dev-ec2imds.device
After=dev-ec2imds.device 
